# packaging/linux/systemd/boxpilotd.service
[Unit]
Description=BoxPilot privileged helper (D-Bus activated)
Documentation=https://github.com/BeyondCortex/BoxPilot
After=dbus.service
Requires=dbus.service

[Service]
Type=dbus
BusName=app.boxpilot.Helper
ExecStart=/usr/lib/boxpilot/boxpilotd
User=root
Restart=on-failure
RestartSec=2s
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
RuntimeDirectory=boxpilot
ReadWritePaths=/etc/boxpilot /var/lib/boxpilot /var/cache/boxpilot /etc/systemd/system /etc/polkit-1/rules.d /run/systemd /run/dbus
CapabilityBoundingSet=CAP_SYS_ADMIN CAP_NET_ADMIN CAP_DAC_OVERRIDE CAP_CHOWN CAP_FOWNER

# No [Install] section: this unit is D-Bus activated on demand via the
# SystemdService= line in app.boxpilot.Helper.service. It is never enabled
# at boot, so dh_installsystemd's default `enable` is a no-op (no install
# directives means systemd has nothing to symlink).
